# etc/supervisord.d/postfix.ini
[program:postfix]
process_name = master
command=/usr/sbin/postfix -v -v -c /etc/postfix start
startsecs=0
autorestart=false
