# etc/supervisord.d/rsyslog.ini
[program:rsyslog]
command=/usr/sbin/rsyslogd -n
